# hello.service
[Unit]
Description="MY FIRST FILE"

[Service]
ExecStart=/home/knoldus/hello.sh

[Install]
WantedBy=multi.user.target

